# docs/examples/systemd/dffmpeg-coordinator.service
[Unit]
Description=DFFmpeg Coordinator Service
After=network.target

[Service]
Type=simple
User=dffmpeg
Group=dffmpeg
# Adjust this path if you installed elsewhere
WorkingDirectory=/opt/dffmpeg/coordinator
Environment="PATH=/opt/dffmpeg/coordinator/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/dffmpeg/coordinator/bin/dffmpeg-coordinator
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
